# deploy/gastown-web.service
[Unit]
Description=Gas Town Web UI
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/shisui/gt
ExecStart=/home/shisui/go/bin/gt gui --port 8080
Restart=on-failure
RestartSec=5
Environment=HOME=/home/shisui
Environment=GT_ROOT=/home/shisui/gt
Environment=PATH=/home/shisui/go/bin:/home/shisui/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
Alias=gastown-gui.service
WantedBy=default.target
