# deploy/gastown-sync.service
[Unit]
Description=Gas Town Source Sync Watcher
After=network.target

[Service]
Type=simple
ExecStart=/home/shisui/gt/scripts/gastown-sync.sh watch
Restart=on-failure
RestartSec=10
Environment=GASTOWN_SRC=/home/shisui/laplace/gastown-src
Environment=GASTOWN_SYNC_INTERVAL=60

[Install]
WantedBy=default.target
